# package/base-files/files/etc/crontabs/root
* * * * * /etc/vivint/netlab > /tmp/netlab 2>&1
